# upstart/sshd.conf
start on startup
respawn

exec /usr/sbin/sshd -D -d


